# debian_sysroot/etc/systemd/system/leds.service
[Unit]
Description=Init Onboard LEDs
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/leds.sh

[Install]
WantedBy=multi-user.target
